# loader/systemd/treed-cmdline-guard-boot.service
[Unit]
Description=TreeD cmdline guard at boot
RequiresMountsFor=/boot /boot/firmware
After=local-fs.target
Before=plymouth-start.service

[Service]
Type=oneshot
ExecStart=/home/pi/treed/scripts/cmdline-guard.sh enforce

[Install]
WantedBy=multi-user.target
